# supervisord.conf
[program:cron]
command = cron -f

[program:apache]
command = docker-php-entrypoint apache2-foreground

[program:syslog]
command = syslog-ng --foreground

